# deploy/gtpconverter.service
[Unit]
Description=GTP Converter - Demucs Backend API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/gtpconverter/server
ExecStart=/opt/gtpconverter/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/gtpconverter
# Загрузить .env если есть (создайте /opt/gtpconverter/.env)
EnvironmentFile=-/opt/gtpconverter/.env

[Install]
WantedBy=multi-user.target
